# launchagent/com.claude-telegram-rs.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <!-- Keep this in sync with the filename you install to (e.g. com.claude-telegram-rs) -->
    <string>com.claude-telegram-rs</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Preferred: run via start-rust.sh (loads .env, sets PATH, executes Rust binary) -->
        <string>/bin/bash</string>
        <string>/Users/USERNAME/Dev/claude-telegram-bot-ts/launchagent/start-rust.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- Update this path to your bot installation -->
    <string>/Users/USERNAME/Dev/claude-telegram-bot-ts</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Optional: If you prefer, you can keep secrets in .env instead of this plist. -->

        <!-- Required: Your Telegram bot token -->
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>your-bot-token-here</string>

        <!-- Required: Comma-separated Telegram user IDs -->
        <key>TELEGRAM_ALLOWED_USERS</key>
        <string>123456789</string>

        <!-- Recommended: Working directory for Claude -->
        <key>CLAUDE_WORKING_DIR</key>
        <string>/Users/USERNAME/Dev</string>

        <!-- Optional: For voice transcription -->
        <key>OPENAI_API_KEY</key>
        <string>sk-...</string>

        <!-- Optional: Comma-separated paths Claude can access -->
        <key>ALLOWED_PATHS</key>
        <string>/Users/USERNAME/Dev,/Users/USERNAME/Documents</string>

        <!-- Optional: Thinking keywords -->
        <key>THINKING_KEYWORDS</key>
        <string>think,pensa,ragiona</string>
        <key>THINKING_DEEP_KEYWORDS</key>
        <string>ultrathink,think hard,pensa bene</string>

        <!-- Optional: Rate limiting -->
        <key>RATE_LIMIT_ENABLED</key>
        <string>true</string>
        <key>RATE_LIMIT_REQUESTS</key>
        <string>40</string>
        <key>RATE_LIMIT_WINDOW</key>
        <string>60</string>

        <!-- Ensure PATH includes necessary binaries (claude, pdftotext, etc.) -->
        <key>PATH</key>
        <string>/Users/USERNAME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/claude-telegram-rs.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-telegram-rs.err</string>
</dict>
</plist>
